# root/etc/services.d/murmur/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc /usr/bin/murmurd \
  -fg \
  -ini /config/murmur.ini
